# debian/debian/statsd/etc/supervisor/conf.d/statsd.conf
[program:statsd]
command=/usr/share/statsd/scripts/start
user=statsd
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=QUIT